# ansible/configs/patch-promoter.path
[Unit]
Description=Change tcp to rdma

[Path]
PathChanged=/etc/drbd-reactor.d
TriggerLimitIntervalSec=0

[Install]
WantedBy=multi-user.target